# mac/bluetooth-audio-fix/com.user.bluetooth.audio.fix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.bluetooth.audio.fix</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>~/.bluetooth-audio-fix/bluetooth_audio_monitor.sh</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>BLUETOOTH_DEVICE_NAME</key>
        <string></string>
        <key>INPUT_VOLUME</key>
        <string>100</string>
        <key>PREFERRED_OUTPUT_DEVICE</key>
        <string></string>
    </dict>
    
    <!-- 监听蓝牙相关的系统文件变化 -->
    <key>WatchPaths</key>
    <array>
        <string>/Library/Preferences/com.apple.Bluetooth.plist</string>
        <string>/Library/Preferences/Audio</string>
    </array>
    
    <!-- 限制执行频率，避免频繁触发 -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>~/Library/Logs/bluetooth_audio_fix_output.log</string>
    
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/bluetooth_audio_fix_error.log</string>
</dict>
</plist>
